# debian/cron.d
30	*	*	*	*	root	curl localhost:13710/api/update
*	*	*	*	*	root	curl -s localhost:13710/api/ping | fgrep ok || service mosecom-air restart
